# vmidentity/config/vmware-sampled.service.in
[Unit]
Description=VMware Secure Token Sample Service
After=syslog.target network.target lwsmd.service
Requires=lwsmd.service

[Service]
Type=forking
ExecStart=/opt/vmware/sbin/vmware-sampled.sh start
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target
