# vmidentity/config/vmware-stsd.service.in
[Unit]
Description=VMware Secure Token Service
After=syslog.target network.target lwsmd.service
Requires=lwsmd.service

[Service]
Type=forking
ExecStart=/opt/vmware/sbin/vmware-stsd.sh start
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
