# vmidentity/config/vmware-idmd.service.in
[Unit]
Description=VMware Identity Manager Service
After=syslog.target network.target

[Service]
Type=forking
ExecStart=@VMSTS_SBIN_DIR@/vmware-idmd.sh start
PIDFile=/var/run/vmware-idmd.pid

[Install]
WantedBy=multi-user.target
